# etc/supervisor/conf.d/kibana.conf
[program:kibana]
command=/opt/kibana/bin/kibana
stdout_logfile=/var/log/kibana/stdout.log
redirect_stderr=true
